# etc/local.d/openpdu.start
#!/bin/sh

/usr/bin/openpdu initboards
/usr/bin/openpdu initoutlets
